# Smart_Assitant/com.smartexit.assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Unique service identifier -->
    <key>Label</key>
    <string>com.smartexit.assistant</string>

    <!-- Command to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/madhulatha/.pyenv/versions/3.10.13/bin/python</string>
        <string>-m</string>
        <string>streamlit</string>
        <string>run</string>
        <string>/Users/madhulatha/Downloads/Python_project/app.py</string>
        <string>--server.port=8501</string>
        <string>--server.headless=true</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/madhulatha/Downloads/Python_project</string>

    <!-- Start automatically on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log files -->
    <key>StandardOutPath</key>
    <string>/Users/madhulatha/Downloads/Python_project/app.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/madhulatha/Downloads/Python_project/app_error.log</string>

    <!-- Wait 5 seconds before restarting on crash -->
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/madhulatha/.pyenv/versions/3.10.13/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/madhulatha</string>
    </dict>

</dict>
</plist>
